# superpeer.service
[Unit]
Description=RightMesh Superpeer
[Service]
ExecStart=/home/ubuntu/superpeer/build/install/Superpeer/bin/Superpeer --headless
User=ubuntu
